# packages/spec/src/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
# user = node
directory = /usr/src/app
command = node index.js
# command = node index.js
